# demo/roles/nginx/templates/supervisor.conf
[program:nginx]
command={{ conda_prefix }}/bin/nginx -p {{ prefix }} -c {{ prefix }}/etc/nginx/nginx.conf
directory={{ prefix }}
