# demo/roles/wps/templates/supervisor.conf.j2
[program:{{ program_name }}]
command=bash run.sh
autostart=true
directory={{ project_location }}
redirect_stderr=true
stdout_logfile=/var/log/supervisor/{{ program_name }}.log
